# supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
stdout_logfile=/var/log/supervisor/nginx.log
stderr_logfile=/var/log/supervisor/nginx.err
autorestart=true
priority=30

[program:php-fpm]
command=/usr/local/bin/php-fpm-run -F
stdout_logfile=/var/log/supervisor/php-fpm.log
stderr_logfile=/var/log/supervisor/php-fpm.err
autorestart=true
priority=20

[program:flask]
command=gunicorn --workers 2 --bind 127.0.0.1:5001 app:app
directory=/app/backend
stdout_logfile=/var/log/supervisor/flask.log
stderr_logfile=/var/log/supervisor/flask.err
autorestart=true
priority=10
